# scripts/com.patchwork.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.patchwork.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jonogompels/.patchwork/watchdog.sh</string>
    </array>

    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Trigger immediately when settings.json changes (e.g. Claude Code update) -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/jonogompels/.claude/settings.json</string>
    </array>

    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/Users/jonogompels/.patchwork/watchdog-stderr.log</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
